# cron/crontab
#crontab
*/5 * * * *	rif *applications/sentinel/modules/start_readers.py