# cron/crontab
#crontab
*/5 * * * *	rif *applications/sentinel/cron/start_readers.py